# systemd_units/chrony-tracking.service
[Unit]
Description=Parse chrony tracking log and emit metrics

[Service]
Type=oneshot
ExecStart=/lib/rdk/chrony_metrics_parse.sh
