# systemd_units/chrony-tracking.path
[Unit]
Description=Watch chrony tracking log for updates

[Path]
PathModified=/opt/logs/chrony/tracking.log

[Install]
WantedBy=multi-user.target
