# systemd_units/myservice.service
[Unit]
Description=This is for Tesing. Do not Merge
After=ntp-sync.target

[Service]
Type=simple
ExecStart=/bin/sh -c 'exit 1'


[Install]                                          
WantedBy=ntp-sync.target
