# systemd_units/NetworkManager.path
[Unit]
Description= NetworkManager Path

[Path]
PathExists=/tmp/bootType
Unit=NetworkManager.service

[Install]
WantedBy=multi-user.target
