# systemd_units/ntp-time-sync.target
[Unit]
Description=Target For NTP Time sync
