# systemd_units/ntp-response-telemetry.service
[Unit]
Description=Handle NTP status log update

[Service]
Type=oneshot
ExecStart=/lib/rdk/update-ntp-response-telemetry.sh
