# systemd_units/ntp-sync.target
[Unit]
Description=Custom target after NTP sync event
