# systemd_units/ntp-telemetry-log.path
[Unit]
Description=NTP Server Response Event
OnFailure=path-fail-notifier@%n.service

[Path]
PathChanged=/tmp/ntp-status-latest.log
Unit=ntp-response-telemetry.service

[Install]
WantedBy=multi-user.target
