# rdktv/systemd_units/vdec-statistics.service
[Unit]
Description=Vdec statistics for telemtry

[Service]
Type=oneshot
ExecStart=/lib/rdk/vdec-statistics.sh

[Install]
WantedBy=multi-user.target
